# .config/systemd/user/nextpaper.service
[Unit]
Description=Changes wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/nextpaper.sh %h/bg/
