# .config/systemd/user/nextpaper.timer
[Unit]
Description=Change Background periodically
Requires=nextpaper.service

[Timer]
AccuracySec=10s
OnUnitActiveSec=1m
